# kovan-failover-node/scripts/parity.service
[Unit]
Description=Parity Daemon
After=network.target

[Service]
EnvironmentFile=/home/USERNAME/parity.conf
ExecStart=/usr/bin/parity $ARGS
User=USERNAME

[Install]
WantedBy=default.target

